# templates/default/supervisord.realtime.conf.erb
[program:realtime]
command = /usr/bin/java -Xmx256m -Duser.timezone=UTC -Dfile.encoding=UTF-8 -Ddruid.realtime.specFile=/var/realtime.spec -classpath /var/druid-<%=@version%>/lib/*:/var/druid-<%=@version%>/config/realtime io.druid.cli.Main server realtime
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/realtime.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/realtime.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:realtime_server]
programs=realtime




